# .config/autostart/my-xmodmap.desktop
[Desktop Entry]
Name=MyXmodmap
Exec=/usr/bin/xmodmap /home/rrc/.Xmodmap
Terminal=false
Type=Application
